# deploy/templates/hstd/clientrotate.service.j2
[Unit]
Description=Sync xray client conns from subscription

[Service]
Type=oneshot
ExecStart=/usr/bin/xrayvpn sub sync
